# gears/quaily-journalist.service.example
[Unit]
Description=Quaily Journalist Service
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/quaily-journalist
ExecStart=%h/quaily-journalist/bin/quaily-journalist serve --config %h/quaily-journalist/config.yaml
Restart=on-failure
# Configuration is read from YAML files only; no env overrides are used.

[Install]
WantedBy=multi-user.target
